# master/supervisord.conf
[program:pollqlog]
command = python /home/dotcloud/current/pollqlog.py
